# f20-redis/supervisord-redis.conf
;------------------------------------------------------------------------------------------------------------------------------
; postgres supervisord config file.
;
; For more information on the config file, please see:
; http://supervisord.org/configuration.html
;
; Note: shell expansion ("~" or "$HOME") is not supported.  Environment
; variables can be expanded using this syntax: "%(ENV_HOME)s".
;------------------------------------------------------------------------------------------------------------------------------
[program:redis-server]
command=/bin/bash -c "/usr/sbin/redis-server"
user=redis
